# Software/deploy/pi/systemd/aris-kiosk.service
[Unit]
Description=ARis Kiosk Launcher
After=systemd-user-sessions.service aris-ui.service
Requires=aris-ui.service

[Service]
Type=simple
User=admin
Environment=HOME=/home/admin
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/admin/.Xauthority"
WorkingDirectory=/home/admin
ExecStart=/home/admin/Studienarbeit_ARis/Software/deploy/pi/kiosk/start_onboard.sh
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
