# Software/deploy/pi/systemd/aris-ui.service
[Unit]
Description=ARis Onboard UI (Svelte preview)
After=network-online.target aris-backend.service
Wants=network-online.target
Requires=aris-backend.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/StudienArbeit_ARis/Software/Onboard_UI
Environment=HOST=127.0.0.1
Environment=PORT=4173
ExecStart=/usr/bin/npm run preview -- --host 127.0.0.1 --port 4173
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
