# Software/deploy/pi/systemd/aris-kiosk.service
[Unit]
Description=ARis Kiosk Launcher
After=systemd-user-sessions.service aris-ui.service
Requires=aris-ui.service

[Service]
Type=simple
User=admin
Environment=HOME=/home/admin
Environment=ARIS_KIOSK_BROWSER=cog
Environment=ARIS_COG_LAUNCHER=cage
Environment=ARIS_DISPLAY_ROTATION=inverted
WorkingDirectory=/home/admin
ExecStart=/home/admin/StudienArbeit_ARis/Software/deploy/pi/kiosk/start_onboard.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
