# Software/deploy/pi/systemd/aris-ui.service
[Unit]
Description=ARis Onboard UI (Svelte preview)
After=network-online.target aris-backend.service
Wants=network-online.target
Requires=aris-backend.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/StudienArbeit_ARis/Software/Onboard_UI
Environment=NODE_ENV=development
ExecStart=/home/admin/.nvm/versions/node/v20.20.2/bin/node /home/admin/StudienArbeit_ARis/Software/Onboard_UI/node_modules/vite/bin/vite.js preview --host 0.0.0.0 --port 4173
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
